# home/sufojareth/.config/systemd/user/dotfiles-sync.service
[Unit]
Description=Auto backup dotfiles and scripts to GitHub

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-sync.sh
